# packaging/systemd/media-server.service
[Unit]
Description=Media Server
Documentation=https://github.com/dog4ik/media-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/media-server
WorkingDirectory=/var/lib/media-server

User=media-server
Group=media-server

StateDirectory=media-server
CacheDirectory=media-server
LogsDirectory=media-server
# Make freshly created files group-r/w
UMask=0002

# Optional overrides and secrets (tokens, RUST_LOG, OTLP endpoint, ...).
# Shipped empty; the leading "-" makes it optional. See media-server.env.
EnvironmentFile=-/etc/media-server/media-server.env
Environment=RUST_LOG=media_server=info

Restart=always
RestartSec=5
TimeoutStopSec=30
KillMode=process

StandardOutput=journal
StandardError=journal

# Sandboxing
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
