# packaging/systemd/media-server.service
[Unit]
Description=Media Server
Documentation=https://github.com/dog4ik/media-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/media-server
WorkingDirectory=/usr/share/media-server

# Run as a transient, isolated system user. systemd creates and owns the
# directories below and exports their paths to the process as $STATE_DIRECTORY,
# $CONFIGURATION_DIRECTORY, $CACHE_DIRECTORY and $LOGS_DIRECTORY, which the binary
# reads (see AppResources::systemd_dir in src/config.rs). This replaces the old
# manual user-creation / chown setup.
DynamicUser=yes
StateDirectory=media-server
ConfigurationDirectory=media-server
CacheDirectory=media-server
LogsDirectory=media-server

# Optional overrides and secrets (tokens, RUST_LOG, OTLP endpoint, ...).
# Shipped empty; the leading "-" makes it optional. See media-server.env.
EnvironmentFile=-/etc/media-server/media-server.env
Environment=RUST_LOG=media_server=info

Restart=always
RestartSec=5
TimeoutStopSec=30
KillMode=process

StandardOutput=journal
StandardError=journal

# NOTE: media libraries usually live outside the service's own directories. If
# your media is not world-readable, grant the service access via
# SupplementaryGroups= (add a group that can read your media) or filesystem ACLs.

[Install]
WantedBy=multi-user.target
